# files/loki/service/template.service.j2
[Unit]
Description=mount
After=network.target

[Service]
ExecStart={{ start_command }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
